# deploy/debian/prerm
#!/bin/bash
set -euo pipefail

SERVICE_NAME="fuelflux"

case "${1:-}" in
  upgrade|failed-upgrade|remove|purge|deconfigure)
    if systemctl is-active --quiet "${SERVICE_NAME}" 2>/dev/null; then
      systemctl stop "${SERVICE_NAME}" || true
    fi
    if systemctl is-enabled --quiet "${SERVICE_NAME}" 2>/dev/null; then
      systemctl disable "${SERVICE_NAME}" || true
    fi
    ;;
  *)
    echo "Warning: unhandled prerm argument '${1:-}'" >&2
    ;;
esac
